# docker/php/supervisor.d/websockets.conf
[program:laravel-websockets]
command=/usr/local/bin/php /var/www/app/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/www/app/storage/logs/websockets.log
